# systemd/totify.service
[Unit]
Description=Totem Keyboard Notifier

[Service]
Type=oneshot
ExecStart=totify

